# spryker/files/etc/cron.d/yves-zed-cleanup
PATH=/usr/bin:/usr/sbin:/bin:/sbin

# Remove all but last 5 shop releases
*/4 1 * * *	root	cd /data/shop; for ENV in *; do nice rm -rf `find $ENV/releases -mindepth 1 -maxdepth 1 2>/dev/null | sort -r | tail -n+6`; done

# Remove tomcat logs older than 4 days
9 1 * * *	root	nice rm -rf `find /data/logs/*/tomcat/ -maxdepth 1 -mindepth 1 -mtime +4 2>/dev/null`

# Remove jenkins builds older than 7 days
10 1 * * *	root	nice rm -rf `find /data/shop/*/shared/data/jenkins/jobs/*/builds/ -maxdepth 1 -mindepth 1 -ctime +7 2>/dev/null`
